# rules/windows/process_creation/proc_creation_win_win10_sched_task_0day.yml
title: Windows 10 Scheduled Task SandboxEscaper 0-day
id: 931b6802-d6a6-4267-9ffa-526f57f22aaf
status: test
description: Detects Task Scheduler .job import arbitrary DACL write\par
author: Olaf Hartong
references:
  - https://github.com/SandboxEscaper/polarbearrepo/tree/master/bearlpe
date: 2019/05/22
modified: 2022/05/13
logsource:
  category: process_creation
  product: windows
detection:
  selection_img:
    Image|endswith: '\schtasks.exe'
    OriginalFileName: 'schtasks.exe'
  selection_cli:
    CommandLine|contains|all:
      - '/change'
      - '/TN'
      - '/RU'
      - '/RP'
  condition: all of selection*
falsepositives:
  - Unknown
level: high
tags:
  - attack.privilege_escalation
  - attack.t1053.005
  - car.2013-08-001
